# deploy/wpa-mcp.service
[Unit]
Description=wpa-mcp WiFi control MCP server (Docker + netns phy)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Edit these to match your host. WIFI_INTERFACE in particular must be set
# to the wireless interface you want the container to own.
Environment=WIFI_INTERFACE=wlan0
Environment=WPA_MCP_IMAGE=wpa-mcp:latest
Environment=PORT=3000
Environment=WPA_DEBUG_LEVEL=2
Environment=WPA_MCP_VOLUME=wpa-mcp-data

# Optional: configure the in-container browser's locale and timezone for
# i18n captive-portal testing. Both are read at container start; changing
# them requires a service restart. See
# docs/design/14_Browser_Locale_Timezone_Design.md.
#Environment=WPA_MCP_BROWSER_LANG=pt-PT
#Environment=WPA_MCP_BROWSER_TZ=Europe/Lisbon

ExecStart=/usr/local/sbin/wpa-mcp-start
ExecStop=/usr/bin/docker rm -f wpa-mcp

[Install]
WantedBy=multi-user.target
